# deployments/systemd/blazelog-server.service
[Unit]
Description=BlazeLog Central Server
Documentation=https://github.com/good-yellow-bee/blazelog
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=blazelog
Group=blazelog

# Binary and config
ExecStart=/usr/local/bin/blazelog-server -c /etc/blazelog/server.yaml
ExecReload=/bin/kill -HUP $MAINPID

# Environment
EnvironmentFile=-/etc/blazelog/server.env

# Working directory
WorkingDirectory=/var/lib/blazelog

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
SystemCallArchitectures=native

# Read-write directories
ReadWritePaths=/var/lib/blazelog /var/log/blazelog

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=blazelog-server

[Install]
WantedBy=multi-user.target
